# ilc-audio-recorder.service
[Unit]
Description=ILC Audio Recorder Web Service
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ilc-audio-recorder
Environment="RECORDINGS_DIR=/home/pi/recordings"
ExecStart=/usr/bin/python3 /home/pi/ilc-audio-recorder/app.py

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/var/log/ilc-audio-recorder/app.log
StandardError=append:/var/log/ilc-audio-recorder/error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
